# db.cron
0 0 * * * /app/backup.sh
0 0 * * * /app/remote_backup.sh
